# contrib/echoflow.service
[Unit]
Description=EchoFlow — voice-to-text daemon
Documentation=https://github.com/WaffleWaldo/EchoFlow
After=graphical-session.target pipewire.service

[Service]
Type=simple
ExecStart=%h/.local/bin/echoflow daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
